# systemd/edo-start-up.service
[Unit]
Description=Run Edo startup script after boot delay
After=network.target

[Service]
Type=oneshot
ExecStart=/bin/bash -lc "/bin/sleep 30; /home/edo/rpi/sh/start-up.sh"

[Install]
WantedBy=multi-user.target
