# systemd/edo-reboot.timer
[Unit]
Description=Weekly Edo reboot (Fridays at 02:00)

[Timer]
OnCalendar=Fri *-*-* 02:00:00
Persistent=true
Unit=edo-reboot.service

[Install]
WantedBy=timers.target
