# systemd/edo-reboot.service
[Unit]
Description=Weekly Edo reboot with notification

[Service]
Type=oneshot
ExecStart=/home/edo/rpi/sh/reboot.sh
