# systemd/edo-rpi.service
[Unit]
Description=Run Edo hourly rpi maintenance script

[Service]
Type=oneshot
User=edo
ExecStart=/home/edo/rpi/rpi.sh
