# systemd/edo-reboot.service
[Unit]
Description=Weekly Edo reboot with notification

[Service]
Type=oneshot
User=edo
ExecStart=/home/edo/sh/reboot.sh
