# systemd/edo-update.service
[Unit]
Description=Run Edo daily system update

[Service]
Type=oneshot
ExecStart=/home/edo/rpi/sh/update.sh
