# systemd/edo-rpi.timer
[Unit]
Description=Hourly Edo rpi maintenance

[Timer]
OnCalendar=hourly
Persistent=true
Unit=edo-rpi.service

[Install]
WantedBy=timers.target
