# systemd/edo-free-space.service
[Unit]
Description=Daily Edo free space check

[Service]
Type=oneshot
ExecStart=/home/edo/rpi/sh/free_space.sh
