# systemd/edo-rpi.service
[Unit]
Description=Run Edo hourly rpi maintenance script

[Service]
Type=oneshot
ExecStart=/home/edo/rpi/sh/rpi.sh
